# api/crontab
* * * * * echo "cron"
* * * * * node cron/bin/updateListingStatuses.js
